# src/bare_metal/systemd/mmwave.service
[Unit]
Description=Mesh sensor publisher: Apollo MSR-2 (mmWave)
After=network-online.target nats-broker.service hostapd.service
Wants=network-online.target
Requires=nats-broker.service

[Service]
Type=simple
User=mesh
Group=mesh
EnvironmentFile=/etc/default/iot-mesh
WorkingDirectory=/opt/iot-mesh-cluster
ExecStart=/opt/iot-mesh-cluster/.venv/bin/python -m bare_metal.sensors.mmwave
Restart=on-failure
RestartSec=5
ProtectSystem=strict
ProtectHome=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
